# SBERUniversity/Final-project-ETL/main.cron
0 1 * * * ./main.py > rep.txt && ./teleg.py